# oom-watch/systemd/oom-watch.service
[Unit]
Description=OOM-Protect monitoring daemon (oom-watch)
Documentation=file:/usr/local/share/oom-watch/manual.md https://github.com/milos85vasic/OOM-Protect
# Run early but after the journal is up; we read it on incident.
After=systemd-journald.service
Wants=systemd-journald.service
# Don't compete with system early-boot; this daemon is forensic, not critical.
DefaultDependencies=yes
# StartLimit* MUST live in [Unit] on systemd >= 230. Older versions
# tolerated them under [Service]; systemd 258 (ALT 11) emits
# "Unknown key 'StartLimitIntervalSec' in section [Service], ignoring."
# Three failures within 60 s and the unit gives up — operator must
# investigate (run `systemctl reset-failed oom-watch.service` to clear).
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/local/sbin/oomwatch -config /etc/oom-watch/config.json

# Restart policy: stay up. If we crash three times in a minute, give up so an
# operator notices rather than masking a real bug.
Restart=on-failure
RestartSec=5s

# Permissions: oom-watch needs to read /proc, /sys/fs/cgroup, and to spawn
# atop. atop wants to read /proc/<pid>/stat for every process and /proc/kallsyms,
# which on most distros requires root. We therefore run as root, but lock the
# rest of the unit down hard.
User=root
Group=root

# Hardening — we are root but we do not need to write anywhere except the
# report and state directories.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/log/oom-watch /var/lib/oom-watch
ProtectHome=read-only
PrivateTmp=yes
# ProtectControlGroups: only 'yes' or 'no' are universally accepted.
# Newer systemd added 'private' / 'strict' but older releases (and ALT
# Linux 11's systemd 258) reject anything other than the boolean form.
# 'yes' makes /sys/fs/cgroup read-only inside the unit, which is what
# we need — the daemon reads cgroup counters but never writes them.
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictRealtime=yes
LockPersonality=yes
RestrictSUIDSGID=yes

# NOT enabled (atop compatibility):
#
#   RestrictNamespaces=yes        — atop reads /proc/<pid>/ns/* on some
#                                   hosts to attribute namespaces; this
#                                   directive can block it.
#   SystemCallArchitectures=native — couples to the seccomp filter below.
#   SystemCallFilter=...          — atop occasionally invokes syscalls
#                                   outside @system-service (bpf,
#                                   perf_event_open, quotactl on some
#                                   builds). systemd kills with SIGSYS
#                                   when a non-allowed syscall is made,
#                                   producing 'atop exit -1: ' (signal
#                                   kill, empty stderr) in the daemon log.
#   CapabilityBoundingSet=...     — atop running as root benefits from
#                                   the full capability set for /proc
#                                   reads across containers and namespaces.
#
# These are defence-in-depth on a process that is already root and
# already cgroup-bounded; the security value is small and the
# fragility cost is large. The cgroup MemoryMax/TasksMax and the
# OOMScoreAdjust/ProtectSystem directives above are the load-bearing
# protections. If you re-enable any of these, run
#   make oomwatch-deploy
# and watch journalctl -u oom-watch.service for 'atop exit -1' (signal
# kill — sandbox blocked something) before deploying widely.

# Prevent oom-watch from being itself OOM-killed before it can write a report.
OOMScoreAdjust=-500

# Resource bounds (the daemon is light by design — these prevent runaway).
MemoryMax=128M
MemoryHigh=64M
TasksMax=64

[Install]
WantedBy=multi-user.target
