# oom-watch/systemd/oom-watch.service
[Unit]
Description=OOM-Protect monitoring daemon (oom-watch)
Documentation=file:/usr/local/share/oom-watch/manual.md https://github.com/milos85vasic/OOM-Protect
# Run early but after the journal is up; we read it on incident.
After=systemd-journald.service
Wants=systemd-journald.service
# Don't compete with system early-boot; this daemon is forensic, not critical.
DefaultDependencies=yes
# StartLimit* MUST live in [Unit] on systemd >= 230. Older versions
# tolerated them under [Service]; systemd 258 (ALT 11) emits
# "Unknown key 'StartLimitIntervalSec' in section [Service], ignoring."
# Three failures within 60 s and the unit gives up — operator must
# investigate (run `systemctl reset-failed oom-watch.service` to clear).
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/local/sbin/oomwatch -config /etc/oom-watch/config.json

# Restart policy: stay up. If we crash three times in a minute, give up so an
# operator notices rather than masking a real bug.
Restart=on-failure
RestartSec=5s

# Permissions: oom-watch needs to read /proc, /sys/fs/cgroup, and to spawn
# atop. atop wants to read /proc/<pid>/stat for every process and /proc/kallsyms,
# which on most distros requires root. We therefore run as root, but lock the
# rest of the unit down hard.
User=root
Group=root

# Hardening — we are root but we do not need to write anywhere except the
# report and state directories.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/log/oom-watch /var/lib/oom-watch
ProtectHome=read-only
PrivateTmp=yes
# ProtectControlGroups: only 'yes' or 'no' are universally accepted.
# Newer systemd added 'private' / 'strict' but older releases (and ALT
# Linux 11's systemd 258) reject anything other than the boolean form.
# 'yes' makes /sys/fs/cgroup read-only inside the unit, which is what
# we need — the daemon reads cgroup counters but never writes them.
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_DAC_READ_SEARCH

# Prevent oom-watch from being itself OOM-killed before it can write a report.
OOMScoreAdjust=-500

# Resource bounds (the daemon is light by design — these prevent runaway).
MemoryMax=128M
MemoryHigh=64M
TasksMax=64

[Install]
WantedBy=multi-user.target
